# userspace/systemd/user/mobi.phosh.OSK.service.d/shutdown.conf
# Stop the on-screen keyboard cleanly on shutdown / reboot.
#
# Without this drop-in, mobi.phosh.OSK.service keeps restarting after phoc
# (the Wayland compositor) has already exited at shutdown. Each restart fails
# immediately with "cannot open display: :1", and systemd retries up to the
# StartLimitBurst before giving up. This pads every reboot with ~5-7 seconds
# of restart spam and shows up to the user as jobuser@1000.service/stop
# hanging while user@1000 refuses to terminate.
#
# Same root cause is hit on other Mobian/Phosh devices; this drop-in is a
# generic Phosh fix, not dipper-specific.

[Unit]
Conflicts=shutdown.target
Before=shutdown.target
